# src/delay_collector/delay-collector.service
[Unit]
Description=Busurbano Delay Collector Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=app
Group=app
WorkingDirectory=/opt/delay_collector
EnvironmentFile=/opt/delay_collector/.env

# Restart policy - always restart if it crashes
Restart=always
RestartSec=10

# Start the Python script
ExecStart=/usr/bin/python3 /opt/busurbano/delay_collector/main.py

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=delay-collector

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/busurbano/delay_collector

[Install]
WantedBy=multi-user.target
